# deploy/hetzner/slt-studio.service
[Unit]
Description=Sweet Little Trauma Studio (slt-studio-v2)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/slt-studio-v2
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=SLT_STATIC_DIR=./dist
ExecStart=/usr/bin/node server/api-proxy.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
